# singlehost/components/rabbitmq/config/cloudify-rabbitmq.service
[Unit]
Description=RabbitMQ Service

[Service]
TimeoutStartSec=0
Restart=on-failure
EnvironmentFile=-/etc/sysconfig/cloudify-rabbitmq
ExecStart=rabbitmq-server

[Install]
WantedBy=multi-user.target